# build/live-build/config/includes.chroot/etc/xdg/autostart/calamares.desktop
[Desktop Entry]
Type=Application
Name=JonasOS Installer
Exec=pkexec calamares
Icon=calamares
Terminal=false
Categories=System;
X-GNOME-Autostart-enabled=true
